# init/rtc-ds3231-sync.service
[Unit]
Description=Synchronize system clock from RTC

[Service]
ExecStart=/sbin/hwclock --hctosys --utc -f /dev/rtc1
Type=oneshot
